# cluster/config/systemd/logrotate.timer
[Unit]
Description=Log Rotation

[Timer]
OnCalendar=*:0/15
Persistent=true
